# tools/builder/assets/attach-gadget.service
[Unit]
Description=Attach gadget to UDC (bind) and fix endpoint perms
After=setup-gadget.service ffs_registrar.service
Requires=ffs_registrar.service
# BindsTo ensures this service stops when ffs_registrar stops
BindsTo=ffs_registrar.service
# PartOf ensures this service restarts when ffs_registrar restarts
PartOf=ffs_registrar.service

[Service]
Type=exec
ExecStartPre=/usr/bin/udevadm settle --timeout=5
ExecStart=/usr/local/bin/attach-gadget.sh
RemainAfterExit=no
# Allow time for endpoint creation (script has 60s timeout per endpoint)
TimeoutStartSec=90

Restart=always
RestartSec=1
StartLimitIntervalSec=0

StandardOutput=journal+console
StandardError=journal+console

[Install]
WantedBy=multi-user.target